# masterfiles_default/etc/cron.d/cfengine
# git pull cfengine
*/3 * * * * root cd /var/cfengine/masterfiles/ && /usr/bin/git pull --ff-only origin master 
